# .config/autostart/keyboard.desktop
[Desktop Entry]
Type=Application
Exec="/home/feofan/.keyboard.sh"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Keyboard config
